# etc/init.d/rc.local
#!/bin/sh

cd /usr/local/watchman-alerting/Alert_Zone
/usr/bin/supervisor zone-server.js

cd /