# etc/profile.d/watchman-alerting.csh
#!/bin/csh
setenv PATH /usr/local/watchman-alerting/bin':'$PATH