# etc/profile.d/watchman-alerting.sh
#!/bin/sh
PATH=/usr/local/watchman-alerting/bin':'$PATH
export PATH